# php-alpine/local/fsroot/etc/service/php-fpm/run
#!/bin/sh

mkdir -p /var/log/php

php-fpm -F
