# php-alpine/debug/fsroot/etc/service/php-fpm/run
#!/bin/sh


main() {
  mkdir -p /var/log/php
  php-fpm -F
}

main >> /var/log/svc/php-fpm.log 2>&1 &
